# sys-power/phc-k8/files/init
#!/sbin/runscript
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
    # start() writes to the "/proc" pseudo-filesystem.
	need localmount
}

start() {
    # If "/etc/conf.d/phc-k8" has uncommented ${PHC_K8_VIDS}, revolt *ALL* CPUs.
    if [ -n  "${PHC_K8_VIDS}" ]; then
	    ebegin "Revolting processors to new voltage IDs \"${PHC_K8_VIDS}\""

        local filename
        for   filename in /sys/devices/system/cpu/cpu*/cpufreq/phc_vids; do
            echo "${PHC_K8_VIDS}" > "${filename}"
        done

        eend
    fi
}

stop() {
    # If "/etc/conf.d/phc-k8" has uncommented ${PHC_K8_VIDS}, *ALL* CPUs were
    # revolted by start(). In such case, revert back to the default vids here.
    if [ -n  "${PHC_K8_VIDS}" ]; then
        ebegin "Reverting processors to default voltage IDs"

        local dirname
        for   dirname in /sys/devices/system/cpu/cpu*/cpufreq; do
            cat "${dirname}/phc_default_vids" > "${dirname}/phc_vids"
        done

        eend
    fi
}
